# wb-back/cron.d/laravel-schedule
* * * * * cd /var/www/wb-back && php artisan schedule:run --no-ansi >> /var/www/wb-back/storage/logs/cron.log 2>&1
